# systemd/healthlink-frontend.service
[Unit]
Description=HealthLink Pro Next.js Frontend
Documentation=https://github.com/yourusername/Healthlink_RPC
After=network.target healthlink-middleware.service
Wants=healthlink-middleware.service

[Service]
Type=simple
User=healthlink
Group=healthlink
WorkingDirectory=/opt/healthlink/frontend

# Environment variables
Environment="NODE_ENV=production"
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="PORT=9002"
Environment="NEXT_TELEMETRY_DISABLED=1"
EnvironmentFile=/opt/healthlink/.env

# Node.js memory optimization (low-spec)
Environment="NODE_OPTIONS=--max-old-space-size=256"

# Start command (standalone mode)
ExecStart=/usr/bin/node server.js

# Alternative if not using standalone build
# ExecStart=/usr/bin/npm start

# Restart policy
Restart=on-failure
RestartSec=10s
StartLimitBurst=5
StartLimitIntervalSec=300

# Resource limits
MemoryMax=384M
MemoryHigh=320M
CPUQuota=50%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=healthlink-frontend

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/healthlink/frontend/.next
ReadWritePaths=/opt/healthlink/frontend/public

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
