# systemd/healthlink-middleware.service
[Unit]
Description=HealthLink Pro Middleware API
Documentation=https://github.com/yourusername/Healthlink_RPC
After=network.target docker.service
Requires=docker.service
Wants=healthlink-fabric.service

[Service]
Type=simple
User=healthlink
Group=healthlink
WorkingDirectory=/opt/healthlink/middleware-api

# Environment variables
Environment="NODE_ENV=production"
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/opt/healthlink/.env

# Node.js memory optimization flags (low-spec)
Environment="NODE_OPTIONS=--max-old-space-size=512 --optimize_for_size --gc_interval=100"

# Start command
ExecStart=/usr/bin/node src/server.js

# Restart policy
Restart=on-failure
RestartSec=10s
StartLimitBurst=5
StartLimitIntervalSec=300

# Resource limits (optional, for extra safety)
MemoryMax=768M
MemoryHigh=640M
CPUQuota=75%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=healthlink-middleware

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/healthlink/middleware-api/data
ReadWritePaths=/opt/healthlink/middleware-api/logs

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
